# pi-client/systemd/smarthome-ble.service
[Unit]
Description=SwitchBot CO2センサー BLE スキャナー
After=network-online.target bluetooth.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smarthome-pi-client
EnvironmentFile=/home/pi/smarthome-pi-client/.env

# uv で仮想環境を使う場合
ExecStart=/home/pi/.local/bin/uv run python ble_scanner.py

Restart=on-failure
RestartSec=10

# BLE アクセスに必要
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN

[Install]
WantedBy=multi-user.target
